# root/etc/s6-overlay/s6-rc.d/acng/run
#! /usr/bin/execlineb -P

with-contenv

/usr/sbin/apt-cacher-ng -v -c /etc/apt-cacher-ng foreground=1
