# root/etc/cron.daily/acng-list-update
#! /command/with-contenv /bin/sh

${SCRIPTS_PATH}/update.sh
